# Day5-systemctl-storage-LVM/systemctl/assignment/etc/systemd/system/ping-google.timer
[Unit]
Description=Run ping-google.service every 5 seconds

[Timer]
OnBootSec=5s
Unit=ping.google.service

[Install]
WantedBy=timers.target
