# Day5-systemctl-storage-LVM/systemctl/assignment/etc/systemd/system/ping.google.service
[Unit]
Description=Ping google once


[Service]
Type=oneshot
ExecStart=/usr/local/bin/ping.google.sh
